# applications/.config/autostart/spotify.desktop
[Desktop Entry]
Type=Application
Exec=blockify
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Application Name
Comment=Description
